# city-pop.service
[Unit]
Description=City Pop Pi-Tail deck
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/city-pop
ExecStart=/opt/city-pop/.venv/bin/python /opt/city-pop/app.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
